# kilo_complete_fixes/systemd_services/android_eyes_bridge.service
[Unit]
Description=Kilo Android Eyes Display Bridge
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/kilo/bridges
Environment=PATH=/opt/kilo/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/kilo/venv/bin/python /opt/kilo/bridges/android_eyes_bridge.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target