# kilo personality/systemd/kilo-chat@.service
[Unit]
Description=Kilo Chat HTTP for user %i
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/kilo/personality
ExecStart=/usr/bin/env python3 /opt/kilo/personality/kilo_chat_http.py
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
